# docker/s6-overlay/s6-rc.d/setup_datadir/up
foreground { echo "Calling setup_datadir.sh..." }
/etc/s6-overlay/scripts/setup_datadir.sh
