# docker/s6-overlay/s6-rc.d/pnpm_updates/up
foreground { echo "Calling pnpm_update.sh..." }
/etc/s6-overlay/scripts/pnpm_update.sh
